# deploy/systemd/radar-api.service
[Unit]
Description=Alkira Radar FastAPI
After=network.target

[Service]
Type=simple
User=radar
WorkingDirectory=/opt/radar/api
EnvironmentFile=/opt/radar/.env
ExecStart=/opt/radar/api/.venv/bin/uvicorn radar.api:app --host 127.0.0.1 --port 8601 --workers 1
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
